# sources/meta-kria/recipes-bsp/mac-config/files/update-mac-config.service
[Unit]
Description=update-mac-config

[Service]
ExecStart=/usr/bin/update-mac-config.sh

[Install]
WantedBy=multi-user.target
